# debian/wpasupplicant.upstart
start on started dbus
stop on stopping dbus

respawn

script
	if [ -r /system/etc/ubuntu/device_environment ]; then
		. /system/etc/ubuntu/device_environment
		export ANDROID_HAL_WIFI_DRIVER_ID
	fi

	exec /sbin/wpa_supplicant -u -s -O /run/wpa_supplicant
end script
